# client/supervisor.conf
[program:mediasoup-poc-client]
process_name=%(program_name)s
command=pnpm dev --host
directory=/home/ubuntu/workspace/mediasoup-poc/client
user=ubuntu
autostart=true
autorestart=true
startretries=3
stopwaitsecs=10
stopsignal=SIGTERM
stdout_logfile=/home/ubuntu/workspace/mediasoup-poc/client/logs/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/home/ubuntu/workspace/mediasoup-poc/client/logs/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
environment=PATH="/home/ubuntu/.nvm/versions/node/v22.23.1/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
